# launchd/com.masond84.sag-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.masond84.sag-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@SAG_ROOT@/scripts/launchd-run.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@SAG_ROOT@</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>@LOG_DIR@/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>@LOG_DIR@/stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
